# rootfs/etc/services.d/erbium/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting Erbium DNS/DHCP server..."

CONFIG_FILE="/data/erbium.conf"
DHCP_START=$(bashio::config 'dhcp_range_start')
DHCP_END=$(bashio::config 'dhcp_range_end')
GATEWAY=$(bashio::config 'gateway')

# Use first 3 octets + .0/24
NETWORK=$(echo "${DHCP_START}" | cut -d. -f1-3).0/24

bashio::log.info "Generating Erbium configuration..."

cat > "${CONFIG_FILE}" <<'YAML_EOF'
# Erbium configuration generated by Home Assistant addon
YAML_EOF

echo "addresses: [${NETWORK}]" >> "${CONFIG_FILE}"

# DNS forwarding
cat >> "${CONFIG_FILE}" <<'YAML_EOF'

dns-routes:
YAML_EOF

DNS_SERVER_COUNT=$(bashio::config 'dns_servers | length')
if [ "${DNS_SERVER_COUNT}" -gt 0 ]; then
    FIRST_DNS=$(bashio::config 'dns_servers[0]')
    cat >> "${CONFIG_FILE}" <<YAML_EOF
  - domain-suffixes: [""]
    type: forward
    dns-servers: [${FIRST_DNS}]
YAML_EOF
fi

# Advertise our own address as DNS server to clients
cat >> "${CONFIG_FILE}" <<'YAML_EOF'

dns-servers: [$self4, $self6]
YAML_EOF

DOMAIN=$(bashio::config 'domain')
if [ -n "${DOMAIN}" ] && [ "${DOMAIN}" != "local" ]; then
    cat >> "${CONFIG_FILE}" <<YAML_EOF

dns-search: [${DOMAIN}]
YAML_EOF
fi

if bashio::config.true 'dhcp_enabled'; then
    bashio::log.info "Configuring DHCP..."
    cat >> "${CONFIG_FILE}" <<YAML_EOF

dhcp-policies:
  - policies:
      - match-subnet: ${NETWORK}
        apply-range:
          start: ${DHCP_START}
          end: ${DHCP_END}
        apply-routes:
          - prefix: 0.0.0.0/0
            next-hop: ${GATEWAY}
YAML_EOF
fi

if bashio::config.true 'router_advertisements_enabled'; then
    bashio::log.info "Configuring Router Advertisements..."
    # TODO: Router Advertisements need per-interface config.
    # Placeholder until we figure out how to detect/configure interfaces.
    cat >> "${CONFIG_FILE}" <<'YAML_EOF'

# router-advertisements:
#   eth0:
#     hop-limit: 64
#     managed: false
#     lifetime: 1h
YAML_EOF
fi

if bashio::config.has_value 'custom_config'; then
    bashio::log.info "Appending custom configuration..."
    CUSTOM_CONFIG=$(bashio::config 'custom_config')
    cat >> "${CONFIG_FILE}" <<YAML_EOF

${CUSTOM_CONFIG}
YAML_EOF
fi

bashio::log.info "Configuration written to ${CONFIG_FILE}"
bashio::log.debug "Configuration contents:"
if bashio::config.true 'log_level' 'debug'; then
    cat "${CONFIG_FILE}"
fi

exec /usr/local/bin/erbium "${CONFIG_FILE}"
