# rootfs/etc/services.d/erbium/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Erbium service
# Service management approach inspired by ha-addon-dnsmasq-dhcp by f18m
# https://github.com/f18m/ha-addon-dnsmasq-dhcp
# ==============================================================================

bashio::log.info "Starting Erbium DNS/DHCP server..."

# Get configuration options
CONFIG_FILE="/data/erbium.conf"
DHCP_START=$(bashio::config 'dhcp_range_start')
DHCP_END=$(bashio::config 'dhcp_range_end')
GATEWAY=$(bashio::config 'gateway')

# Calculate network from DHCP range
# Simple approach: use first 3 octets + .0/24
NETWORK=$(echo "${DHCP_START}" | cut -d. -f1-3).0/24

# Generate configuration file in YAML format (Erbium uses YAML despite .conf extension)
bashio::log.info "Generating Erbium configuration..."

cat > "${CONFIG_FILE}" <<'YAML_EOF'
## Erbium configuration generated by Home Assistant addon
## See https://github.com/isomer/erbium for documentation

## The addresses to hand out via DHCP and Router Advertisements
YAML_EOF

echo "addresses: [${NETWORK}]" >> "${CONFIG_FILE}"

# DNS Routes configuration
bashio::log.info "Configuring DNS forwarding..."
cat >> "${CONFIG_FILE}" <<'YAML_EOF'

## DNS routing configuration
dns-routes:
YAML_EOF

# Add DNS servers as forwarding targets
DNS_SERVER_COUNT=$(bashio::config 'dns_servers | length')
if [ "${DNS_SERVER_COUNT}" -gt 0 ]; then
    FIRST_DNS=$(bashio::config 'dns_servers[0]')
    cat >> "${CONFIG_FILE}" <<YAML_EOF
  - domain-suffixes: [""]
    type: forward
    dns-servers: [${FIRST_DNS}]
YAML_EOF
fi

# DNS servers for DHCP/RA advertisements
cat >> "${CONFIG_FILE}" <<'YAML_EOF'

## DNS servers to advertise to clients
## $self4 and $self6 refer to the server's own addresses
dns-servers: [$self4, $self6]
YAML_EOF

# DNS search domain
DOMAIN=$(bashio::config 'domain')
if [ -n "${DOMAIN}" ] && [ "${DOMAIN}" != "local" ]; then
    cat >> "${CONFIG_FILE}" <<YAML_EOF

## DNS search domain
dns-search: [${DOMAIN}]
YAML_EOF
fi

# DHCP Policy configuration
if bashio::config.true 'dhcp_enabled'; then
    bashio::log.info "Configuring DHCP..."
    cat >> "${CONFIG_FILE}" <<YAML_EOF

## DHCP configuration
dhcp-policies:
  - policies:
      - match-subnet: ${NETWORK}
        apply-range:
          start: ${DHCP_START}
          end: ${DHCP_END}
        apply-routes:
          - prefix: 0.0.0.0/0
            next-hop: ${GATEWAY}
YAML_EOF
fi

# Router Advertisements
if bashio::config.true 'router_advertisements_enabled'; then
    bashio::log.info "Configuring Router Advertisements..."
    # Note: Router Advertisements are configured per-interface
    # This is a placeholder - actual interface config may need to be added
    cat >> "${CONFIG_FILE}" <<'YAML_EOF'

## Router Advertisements (configured per interface)
## Uncomment and customize for your network interfaces
# router-advertisements:
#   eth0:
#     hop-limit: 64
#     managed: false
#     lifetime: 1h
YAML_EOF
fi

# Append custom configuration if provided
if bashio::config.has_value 'custom_config'; then
    bashio::log.info "Appending custom configuration..."
    CUSTOM_CONFIG=$(bashio::config 'custom_config')
    cat >> "${CONFIG_FILE}" <<YAML_EOF

## Custom configuration from addon options
${CUSTOM_CONFIG}
YAML_EOF
fi

bashio::log.info "Configuration written to ${CONFIG_FILE}"
bashio::log.debug "Configuration contents:"
if bashio::config.true 'log_level' 'debug'; then
    cat "${CONFIG_FILE}"
fi

# Start Erbium
exec /usr/local/bin/erbium "${CONFIG_FILE}"
