# rootfs/etc/services.d/backend/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start backend web service
# ==============================================================================

cd /backend
exec node server.js
